# src/assets/scripts/service-k3s
#!/sbin/openrc-run

# This script is used on Lima to manage k3s; this is necessary as we cannot get
# a reliable connection to the running process.

# shellcheck shell=ksh

depend() {
  after network-online
  want cgroups
}

start_pre() {
  rm -f /tmp/k3s.*
}

supervisor=supervise-daemon
name=k3s
command=/usr/local/bin/k3s
command_args="server --https-listen-port @PORT@ ${K3S_EXEC:-} >>/var/log/k3s 2>&1"
output_log=/var/log/k3s
error_log=/var/log/k3s

pidfile="/var/run/k3s.pid"
respawn_delay=5
respawn_max=0

set -o allexport
if [ -f /etc/environment ]; then source /etc/environment; fi
