# src/assets/scripts/service-k3s
#!/sbin/openrc-run

# This script is used to manage k3s via OpenRC.

# shellcheck shell=ksh

depend() {
  after network-online
  want cgroups
}

start_pre() {
  rm -f /tmp/k3s.*
}

supervisor=supervise-daemon
name=k3s
command=/usr/local/bin/k3s
command_args="server --https-listen-port ${PORT} ${K3S_EXEC:-}"
K3S_LOGFILE="${K3S_LOGFILE:-${LOG_DIR:-/var/log}/${RC_SVCNAME}.log}"
output_log="${K3S_OUTFILE:-${K3S_LOGFILE}}"
error_log="${K3S_ERRFILE:-${K3S_LOGFILE}}"

pidfile="/var/run/k3s.pid"
respawn_delay=5
respawn_max=0

set -o allexport
if [ -f /etc/environment ]; then source /etc/environment; fi
